# cfg/hyberia.service
[Unit]
Description=Hyberia Controller

[Service]
ExecStart=/usr/bin/hyberiactl

[Install]
WantedBy=graphical.target