# cfg/hyberia-audio.service
[Unit]
Description=Hyberia Audio (alsa)
Requires=hyberia-x.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/alsactl restore -d -f /etc/asound.state

[Install]
WantedBy=graphical.target