# cfg/hyberia-x.service
[Unit]
Description=Hyberia X Server

[Service]
ExecStart=/usr/bin/startx
Restart=always

[Install]
WantedBy=graphical.target